# etc/systemd/system/ingenico.service
[Unit]
Description=Ingenico PPP Session

[Service]
Type=forking
ExecStart=/usr/sbin/pppd call ingenico
